# pkg/debian/DEBIAN/postinst
#!/bin/sh

set -ex

systemctl daemon-reload
systemctl enable battery-notify
systemctl start battery-notify
